# transmission/supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=/usr/sbin/openvpn --config /etc/openvpn/client/wormhole.conf
priority=1
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:transmission]
command=/usr/bin/transmission-daemon -f --config-dir /etc/transmission-daemon/
priority=2
autorestart=true
#user=media
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

#[eventlistener:openvpn]
#command=bash -c "printf 'READY\n' && while read line; do kill -SIGQUIT $PPID; done < /dev/stdin"
#events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL

